# deploy/systemd/tecpey-ops-alert-delivery.timer
[Unit]
Description=Retry TecPey Operational Alert Delivery
Documentation=https://github.com/tecpey/Tecpey-Os/blob/main/docs/operations/COMMUNITY_CHALLENGE_SCHEDULER_RUNBOOK.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
AccuracySec=15s
Unit=tecpey-ops-alert-delivery.service

[Install]
WantedBy=timers.target
