# deploy/systemd/tecpey-ops-alert-delivery.service.in
[Unit]
Description=TecPey Operational Alert Delivery
Documentation=https://github.com/tecpey/Tecpey-Os/blob/main/docs/operations/COMMUNITY_CHALLENGE_SCHEDULER_RUNBOOK.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=900
StartLimitBurst=6

[Service]
Type=oneshot
User=@@RUN_USER@@
Group=@@RUN_GROUP@@
WorkingDirectory=@@APP_DIR@@
EnvironmentFile=@@ENV_FILE@@
Environment=NODE_ENV=production
Environment=TECPEY_OPS_STATE_DIR=@@STATE_DIR@@
Environment=NPM_CONFIG_UPDATE_NOTIFIER=false
Environment=NPM_CONFIG_CACHE=/tmp/tecpey-npm-cache
ExecStartPre=@@NPM_BIN@@ run ops:scheduler:env-check
ExecStart=@@NPM_BIN@@ run ops:alerts:deliver
TimeoutStartSec=2min
UMask=0077
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tecpey-ops-alert-delivery
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
RemoveIPC=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
ReadOnlyPaths=@@APP_DIR@@
ReadWritePaths=@@STATE_DIR@@

[Install]
WantedBy=multi-user.target
